# systemd/singbox-manager.service
[Unit]
Description=Sing-Box Manager Backend API
After=network.target

[Service]
Type=simple
User=deck
Group=deck
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 5555
WorkingDirectory=/home/deck/.config/singbox-manager
Restart=on-failure
RestartSec=5

Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/bin:/usr/sbin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
